# src/systemd/opinsys-ktpapi-watcher.path
[Unit]
Description=Watch Opinsys KTPAPI Cmd file changes
Wants=observer.service
After=network.target

[Path]
PathModified=/media/usb1/.opinsys/cmd

[Install]
WantedBy=multi-user.target