# src/systemd/opinsys-ktpapi-watcher.path
[Path]
PathModified=/media/usb1/.opinsys/cmd

[Install]
WantedBy=multi-user.target